# deployment/frontend.service
[Unit]
Description=LayerPainter Frontend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/layerpainter/frontend
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="NEXT_PUBLIC_API_BASE_URL=http://192.168.0.146:8000"
Environment="PORT=3000"
Environment="HOSTNAME=0.0.0.0"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
